# supervisor/program_textbelt.conf
[program:textbelt]
command=node server/app.js
user=nginx
stdout_logfile=syslog
stderr_logfile=syslog
directory=/srv/textbelt
